# hipop/launchd/com.hipop.pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hipop.pull</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>scripts.feishu_pull</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/luke/Downloads/点购工作流/hipop</string>

    <!-- 每 30 分钟跑一次 -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/luke/Downloads/点购工作流/logs/pull.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/luke/Downloads/点购工作流/logs/pull.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/luke/Library/Python/3.9/bin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
    </dict>
</dict>
</plist>
