# hipop/launchd/com.hipop.ziniao.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hipop.ziniao</string>

    <!-- 紫鸟 web_driver 常驻守护：循环巡检 127.0.0.1:18080，断了就 restart。
         占位符路径由 install.sh 按本机仓库路径渲染（见 install.sh render）。 -->
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>hipop.scripts.ziniao_webdriver</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- 开机/登录自启 + 进程崩了 launchd 重新拉起 daemon（daemon 自身常驻）。 -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/ziniao.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/ziniao.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/luke/Library/Python/3.9/bin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
        <key>PYTHONPATH</key>
        <string>__REPO__</string>
    </dict>
</dict>
</plist>
